# .crontab
00 12 * * * cd /Users/hds/Documents/python学习/auto_set_bgImg py bingBGIMG.py
